# launchd/com.framelog.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.framelog.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vedantmodi/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>/Users/vedantmodi/dev/framelog/watcher.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/vedantmodi/Photos/framelog.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/vedantmodi/Photos/framelog.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/vedantmodi/dev/framelog</string>
</dict>
</plist>
